# startup.service
[Unit]
Description=Run custom script at startup

[Service]
User=root
ExecStart=/usr/local/bin/startup.sh

[Install]
WantedBy=default.target
